# gb10-throttle.service
[Unit]
Description=GB10 Bandwidth Throttle Web GUI
After=network.target
Documentation=https://github.com/mcglothi/AIKB

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/gb10-throttle/app.py
WorkingDirectory=/opt/gb10-throttle
Environment=THROTTLE_IFACE=enP7s7
Restart=on-failure
RestartSec=5

# Run as root so tc commands work without extra capabilities setup
User=root

# Clean shutdown
KillMode=mixed
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
